# boxes/x86/suse/root/etc/systemd/system/bundle.automount
Description=Shared KIWI Bundle Directory

[Automount]
Where=/bundle

[Install]
WantedBy=multi-user.target
