# boxes/suse/root/etc/systemd/system/description.automount
Description=Shared KIWI Description Directory

[Automount]
Where=/description

[Install]
WantedBy=multi-user.target
